# overlay/etc/init/node_peterlyons.conf
description "start and stop the peterlyons.com node.js server"
author "Peter Lyons <pete@peterlyons.com>"

start on filesystem and started networking
respawn

chdir /home/plyons/projects/peter_lyons_web_site/code
setuid www-data
setgid www-data
exec /usr/bin/env -i $(cat env.txt | tr -s '\n' ' ') app/server.coffee  >> ../var/log/upstart.log 2>&1
